# deploy/pytalk.service
# PyTalk Daphne Service (Single Worker)
# Place this file at: /etc/systemd/system/pytalk.service
#
# Usage:
#   Start:   systemctl start pytalk
#   Enable:  systemctl enable pytalk
#   Status:  systemctl status pytalk
#   Logs:    journalctl -u pytalk -f

[Unit]
Description=PyTalk ASGI Server (Daphne)
After=network.target redis-server.service postgresql.service
Requires=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pytalk/backend
Environment="PATH=/opt/pytalk/venv/bin"
EnvironmentFile=/opt/pytalk/backend/.env
ExecStart=/opt/pytalk/venv/bin/daphne -b 127.0.0.1 -p 8001 meet.asgi:application
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/pytalk/backend/logs /var/log/pytalk

MemoryMax=400M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
