# deploy/pytalk.service
# PyTalk Daphne Service (Template Unit for Multi-Worker)
# Place this file at: /etc/systemd/system/pytalk@.service
#
# Usage:
#   Start workers:  systemctl start pytalk@8001 pytalk@8002 pytalk@8003 pytalk@8004
#   Enable on boot: systemctl enable pytalk@8001 pytalk@8002 pytalk@8003 pytalk@8004
#   Status:         systemctl status pytalk@8001
#   Logs:           journalctl -u pytalk@8001

[Unit]
Description=PyTalk ASGI Server (Daphne Worker %i)
After=network.target redis.service postgresql.service
Requires=redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pytalk/backend
Environment="PATH=/opt/pytalk/venv/bin"
EnvironmentFile=/opt/pytalk/backend/.env

# Daphne ASGI server — %i is substituted with the port number
ExecStart=/opt/pytalk/venv/bin/daphne \
    -b 127.0.0.1 \
    -p %i \
    --access-log /var/log/pytalk/daphne_access_%i.log \
    meet.asgi:application

Restart=always
RestartSec=5
StandardOutput=append:/var/log/pytalk/daphne_%i.log
StandardError=append:/var/log/pytalk/daphne_error_%i.log

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/pytalk/backend/logs /var/log/pytalk

# Resource limits per worker (4 workers on t3.small with 2GB)
MemoryMax=300M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
