# deploy/pytalk-celery-beat.service
# PyTalk Celery Beat Scheduler
# Place this file at: /etc/systemd/system/pytalk-celery-beat.service

[Unit]
Description=PyTalk Celery Beat Scheduler
After=network.target redis-server.service postgresql.service
Requires=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pytalk/backend
Environment="PATH=/opt/pytalk/venv/bin"
EnvironmentFile=/opt/pytalk/backend/.env

ExecStart=/opt/pytalk/venv/bin/celery \
    -A meet beat \
    --loglevel=info \
    --schedule=/opt/pytalk/backend/celerybeat-schedule

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/pytalk/backend /var/log/pytalk

# Minimal resources — beat is a scheduler, not a worker
MemoryMax=120M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
